# debian/postrm
#!/bin/sh

if [ "$1" = remove ] ; then
    echo remove
    # Remove services
    systemctl disable eblocker-dyndns
fi
