# debian/prerm
#!/bin/sh

if [ "$1" = remove ] ; then
    # Stop service
    systemctl stop eblocker-dyndns
    # Disable service
    systemctl disable eblocker-dyndns
fi

exit 0
